# conf/supervisord.conf
[program:campr]
directory = /srv/campr/
user = root
command = /srv/campr/conf/startup.sh
stdout_logfile = /var/log/campr/stdout.log
stderr_logfile = /var/log/campr/stderr.log

